# etc/frappy.target
[Unit]
Description=Start all SECoP SEC-nodes
Requires=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target
